# opendkim.service.d/user.conf
[Service]
User=opendkim
Group=postfix
RuntimeDirectory=opendkim
RuntimeDirectoryMode=0750
