# nginx.service.d/hardening.conf
[Service]
##############
# Networking #
##############

# PrivateNetwork= service needs access to system network
RestrictAddressFamilies=AF_INET AF_INET6
IPAccounting=yes
# IPAddressAllow= service needs access to all IP addresses
# IPAddressDeny= service needs access to all IP addresses

###############
# File system #
###############
#  Note that the effect of these settings may be undone by privileged processes. In order to
#  set up an effective sandboxed environment for a unit it is thus recommended to combine
#  these settings with either CapabilityBoundingSet=~CAP_SYS_ADMIN or
#  SystemCallFilter=~@mount.

ProtectHome=yes
ProtectSystem=strict
#Additional paths you should allow on your system configuration. Path /run can be allowed only if Apparmor is present and configured to allow only /var/nginx.pid read-write.
#ReadWritePaths=-/run -/var/ngx_pagespeed_cache
ReadWritePaths=-/run/nginx -/var/log/nginx -/usr/share/nginx/ -/var/lib/nginx
PrivateTmp=yes
LogsDirectory=nginx
RuntimeDirectory=nginx

###################
# User separation #
###################

# PrivateUsers= service needs to view other users files (certbot)
# DynamicUser= service has a defined user
NoNewPrivileges=true

###########
# Devices #
###########

PrivateDevices=yes
DevicePolicy=closed
# DeviceAllow= service has PrivateDevices=yes

##########
# Kernel #
##########

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes

########
# Misc #
########

# CapabilityBoundingSet= specified in separate .conf
# AmbientCapabilities= specified in separate .conf
ProtectHostname=yes
ProtectClock=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
#MemoryDenyWriteExecute=yes -conflicts with pagespeed module, possible conflicts with other modules
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
KeyringMode=private

################
# System calls #
################

SystemCallFilter=@system-service
SystemCallFilter=~@resources @privileged
SystemCallArchitectures=native
