# opendmarc.service.d/user.conf
[Service]
User=opendmarc
Group=postfix
RuntimeDirectory=opendmarc
RuntimeDirectoryMode=0750
