# radicale.service.d/hardening.conf
[Service]
IPAddressDeny=any
IPAddressAllow=localhost
